# service/hades.service
[Unit]
Description=Hades manager
After=network-online.target

[Service]
Type=oneshot

User=hades
Group=hades

WorkingDirectory=/opt/hades

ExecStartPre=-/bin/bash /opt/hades/update
ExecStart=/bin/bash /opt/hades/main

Restart=no

[Install]
WantedBy=multi-user.target
