# service/hades.timer
[Unit]
Description=Hourly Hades update

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
Persistent=true
Unit=hades.service

[Install]
WantedBy=timers.target